# keights/resources/usr/lib/systemd/system/keights-templates.target
[Unit]
Description=keights-templates target

[Install]
RequiredBy=keights-kubeadm-init.service
