# fpm/usr/lib/systemd/system/keights-templatize-kube-scheduler.path
[Unit]
Description=keights-templatize-kube-scheduler path watcher

[Path]
PathExists=/run/keights-collector/asg
PathChanged=/run/keights-collector/asg

[Install]
WantedBy=keights-configure.target
