# fpm/usr/lib/systemd/system/keights-manifests.target
[Unit]
Description=keights-manifests target

[Install]
RequiredBy=keights-bootstrap-token-splitter.service
