# fpm/usr/lib/systemd/system/keights-configure.target
[Unit]
Description=keights-configure
Requires=network.target
After=network.target
