# fpm/usr/lib/systemd/system/keights-bootstrap-token-splitter.path
[Unit]
Description=keights-bootstrap-token-splitter path watcher

[Path]
PathExists=/run/kubernetes/bootstrap-token
PathChanged=/run/kubernetes/bootstrap-token

[Install]
WantedBy=keights-configure.target
