# stack/ami/debian/files/etc/systemd/system/docker.service.d/cloudboss.conf
[Service]
ExecStartPre=/bin/rm -rf /var/lib/docker/network
Restart=on-failure
RestartSec=2s
StartLimitInterval=0