# stack/ami/debian/files/etc/profile.d/crictl.sh
export CONTAINER_RUNTIME_ENDPOINT=unix:///run/containerd/containerd.sock
